# autoship-scripts/autoship.service
[Unit]
Description=AutoShip Deployment Handler
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/autoship/main.py
WorkingDirectory=/opt/autoship
EnvironmentFile=/opt/autoship/.env
Restart=on-failure
User=root

# Security Hardening (optional but useful)
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true

ReadWritePaths=/etc/nginx /etc/letsencrypt /var/lib/autoship/certbot \
               /etc/nginx/sites-available /etc/nginx/sites-enabled /var/www/autoship

[Install]
WantedBy=multi-user.target
